# service/camera-display.service
[Unit]
Description=Raspberry Pi Camera Display Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/pi-camera-display
Environment=PYTHONPATH=/opt/pi-camera-display
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/opt/pi-camera-display/venv/bin/python3 /opt/pi-camera-display/main.py
StandardOutput=append:/var/log/camera-display.log
StandardError=append:/var/log/camera-display.log
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
